# docker/root/etc/services.d/dispatcher/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc php /config/www/artisan monitor:run >> /config/log/dispatcher.log
